# rootfs/etc/services.d/nextcloud/finish
#!/bin/bash
# Called when service exits, triggers restart
exit 0
