# rootfs/etc/services.d/nextcloud/run
#!/bin/bash
# Start Nextcloud Desktop Client
exec /usr/bin/nextcloud --background
