# rootfs/etc/cron.d/nextcloud-update
# Run daily at midnight
0 0 * * * root /usr/local/bin/update-nextcloud.sh
